# server/supervisor.ini
[program:recordit]
command=pipenv run uwsgi --socket 127.0.0.1:3031 --wsgi-file wsgi.py --callable app --processes 4 --threads 2 --stats 127.0.0.1:9191
directory=/home/zero/record-it
user=recordit
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true